# supervisor.conf
# Supervisor Configuration for OpplaOneManager
# Place in: /etc/supervisor/conf.d/onemanager.conf
# Then run: sudo supervisorctl reread && sudo supervisorctl update

[program:onemanager-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/onemanager/backend/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=300
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/onemanager/backend/storage/logs/worker.log
stopwaitsecs=3600

[program:onemanager-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while true; do /usr/bin/php /var/www/onemanager/backend/artisan schedule:run --verbose --no-interaction & sleep 60; done"
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/onemanager/backend/storage/logs/scheduler.log
